# conf/cron.daily/pmort.cron
#!/bin/bash

DAYS=1
PMORT_LOG_DIR=/var/log/pmort

find ${PMORT_LOG_DIR} -daystart -mindepth 1 -maxdepth 1 -type d -ctime ${DAYS} -execdir rm -rf "{}" \;
