# conf/cron.daily/pmort.cron
#!/bin/bash

PMORT_LOG_DIR=/var/log/pmort

find ${PMORT_LOG_DIR} -daystart -mindepth 1 -maxdepth 1 -type d -ctime 1 -execdir rm -rf "{}" \;
