# conf/cron.daily/pmort.cron.jinja
#!/bin/bash

DAYS=1
PMORT_OUTPUT_DIR="{{ PARAMETERS['pmort.output_directory'] }}"

find "${PMORT_OUTPUT_DIR}" -daystart -mindepth 1 -maxdepth 1 -type d -ctime ${DAYS} -execdir rm -rf "{}" \;
